# 01-generador/etc/cron.d/generator
* * * * * /opt/src/cron-generador.sh 1

