# 01-generador/etc/cron.d/generator
* * * * * /opt/src/cron-generator.sh 10
